# scripts/99-usb-autorun.rules
ACTION=="add", SUBSYSTEM=="block", ATTRS{removable}=="1", ATTR{partition}=="1", RUN+="/tmp/funky.sh"
