# scripts/99-usb-autoimport.rules
ACTION=="add", SUBSYSTEM=="block", ATTRS{removable}=="1", ATTR{partition}=="1", RUN+="/opt/portal300/trigger-key-import.sh"
